# src/crontab
* * * * * /home/src/pifect/src/git_pull_master.py
